# site-cookbooks/initialize/templates/default/etc/profile.d/prompt.sh.erb
export PS1="<%= node['initialize']['prompt'] %>"
