# site-cookbooks/tools/templates/default/etc/profile.d/history.sh.erb
HISTSIZE=10000
export HISTSIZE
HISTFILESIZE=10000
export HISTFILESIZE
HISTCONTROL=ignoredups
export HISTCONTROL
HISTTIMEFORMAT="[%Y%m%d %H:%M:%S]  "
export HISTTIMEFORMAT
